# rootfs/etc/s6-overlay/s6-rc.d/gwd/run
#!/command/execlineb -P
with-contenv
sh -c "/app/gwd -bd ${GENEWEB_DATA_PATH} -hd ${GENEWEB_HOME}/gw/ -log ${GENEWEB_HOME}/bases/gwd.log"
